# VCK190-Ethernet/2023.2/pl_4x25g_mrmac/Software/Petalinux/components/yocto/layers/meta-system-controller/recipes-utils/labtool-jtag-support/files/xsdb.service
[Unit]
Description=Enable xsdb
After=rpm-autoload.service

[Service]
Type=forking
ExecStart=/usr/bin/xsdb
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
