# VCK190-Ethernet/2023.2/pl_4x25g_mrmac/Software/Petalinux/components/yocto/layers/meta-system-controller/recipes-utils/system-controller-app/files/system_controller.service
[Unit]
Description=system_controller_app
After=rpm-autoload.service xsdb.service

[Service]
Type=forking
ExecStart=/usr/bin/system_controller.sh
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
