# VCK190-Ethernet/2023.2/pl_4x25g_mrmac/Software/Petalinux/components/yocto/layers/meta-system-controller/recipes-apps/rpm-autoload/files/rpm-autoload.service
[Unit]
Description=rpm-autoload

[Service]
Type=oneshot
ExecStart=/usr/bin/rpm-autoload.sh
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
