# VCK190-Ethernet/2023.2/pl_4x25g_mrmac/Software/Petalinux/components/yocto/layers/meta-system-controller/recipes-utils/scweb/files/scwebrun.service
[Unit]
Description=scweb run 
After=system_controller.service

[Service]
Type=forking
ExecStart=/usr/bin/scwebrun.sh
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
